# packages/0confirmation-keeper/zero-keeper.service
[Unit]
Description=0cf Network Keeper
Documentation=https://github.com/0confirmation/0confirmation/tree/master/packages/0confirmation-keeper
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/PATH/TO/0confirmation/
EnvironmentFile=/PATH/TO/0confirmation/packages/0confirmation-keeper/.env
ExecStart=/PATH/TO/bin/node /PATH/TO/0confirmation/packages/0confirmation-keeper keeper.js
StandardOutput=/var/log/0cf-keeper.log
Restart=on-failure

[Install]
WantedBy=multi-user.target
